# cron/crontab
SHELL=/bin/bash
HOME=/Users/lihex
PATH=/usr/bin:/Users/lihex/.dotfiles/bin:/bin

0,30 * * * * say "Time for a break."
0 * * * * proremote
